# .config/autostart/compton.desktop
[Desktop Entry]
Type=Application
Exec=compton
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=compton
Comment=

