# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Exec=dunst
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=dunst
Comment=

